# packaging/systemd/ttsd.service
[Unit]
Description=DAN Glasses Text-to-Speech Service (KittenTTS + gTTS + espeak-ng)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/workspace/dan-glasses
ExecStart=/usr/local/bin/python /home/workspace/dan-glasses/Services/ttsd/ttsd.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target